# rules/sigma/new_scheduled_task.yml
title: New Scheduled Task Creation
id: 1b2c3d4e-5f6a-7b8c-9d0e-1f2a3b4c5d6e
status: experimental
description: Detects the creation of a new scheduled task using the native schtasks.exe utility.
author: Gardiyan Labs
date: 2025-10-17
modified: 2025-10-17
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains: '/create'
  condition: selection
falsepositives:
  - Legitimate administrative scripts or software installation creating tasks.
level: medium
tags:
  - attack.t1053.005
references:
  - https://attack.mitre.org/techniques/T1053/005/
