# rules/sigma/new_scheduled_task.yml
title: New Scheduled Task Creation
id: 1b2c3d4e-5f6a-7b8c-9d0e-1f2a3b4c5d6e
status: experimental
description: Detects creation of a new scheduled task.
author: Your Name <your.email@example.com>
date: 2025/10/17
modified: 2025/10/17
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    EventID: 4698
  condition: selection
falsepositives:
  - Legitimate scheduled task creation
level: medium
tags:
  - attack.persistence
  - attack.t1053
references:
  - https://attack.mitre.org/techniques/T1053/
